# supervisord.conf
[supervisord]
logfile=/dev/null
nodaemon=true

[program:X11]
command=/usr/bin/Xvfb :0 -screen 0 "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s"x24
autorestart=true
user=idea
priority=100

[program:x11vnc]
command=/usr/bin/x11vnc -xkb -noxrecord -noxfixes -noxdamage -display :0 -nopw -wait 5 -shared -permitfiletransfer -tightfilexfer
autorestart=true
user=idea
priority=200

[program:novnc]
command=/root/noVNC/utils/launch.sh --vnc localhost:5900 --listen 8080
autorestart=true

[program:socat]
command=socat tcp-listen:6000,reuseaddr,fork unix:/tmp/.X11-unix/X0
autorestart=true

[program:openbox]
environment=HOME="/home/idea",DISPLAY=":0",USER="idea"
command=/usr/bin/openbox
user=idea
autorestart=true
priority=300

[program:terminal]
command=xfce4-terminal --default-working-directory=/workspace
autorestart=true
user=idea
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
